# ansible-role/templates/piki.target.j2
[Unit]
Description=PiKi Target

[Install]
WantedBy=multi-user.target
